# packaging/app-configuration-backup.cron
50 3 * * * root /usr/clearos/apps/configuration_backup/deploy/backup-now >/dev/null 2>&1
